# infrastructure/somleng/modules/nat_instance/templates/snat.service
[Unit]
Description = SNAT via secondary interface

[Service]
ExecStart = /opt/nat/snat.sh
Type = oneshot

[Install]
WantedBy = multi-user.target
